# Docker/dropbox/periodic/daily/upload
#!/bin/sh
# backup workdir

echo "* Run uploader"
/opt/dropbox_uploader/dropbox_uploader.sh -f /config/.dropbox_uploader upload /workdir /
